# etc/systemd/system/pacman-cache.service
[Unit]
Description=pacman clear package cache

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
