# etc/systemd/system/pacman-cache.timer
[Unit]
Description=pacman clear package cache

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
