# etc/systemd/system/pacman-optimize.service
[Unit]
Description=pacman optimize database

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-optimize
