# etc/systemd/system/pacman-optimize.timer
[Unit]
Description=pacman optimize database

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
